# profile.d/crawlpath.sh
# crawl init env script
#

crawlPath="/home/ihubuser/crawltest_denny/"

export crawlPath


